# configs/supervisor.conf
[supervisord]
nodaemon=true

[program:sqlrely]
command=sqlr-start -id example
redirect_stderr=true
